# .config/systemd/user/feh-wallpaper.timer
[Unit]
Description=Random wallpaper with feh

[Timer]
OnStartupSec=10m
Unit=feh-wallpaper.service

[Install]
WantedBy=default.target
